# services/mysite.socket
[Unit]
Description=gunicorn socket (mysite)

[Socket]
ListenStream=/run/gunicorn/mysite.sock

[Install]
WantedBy=sockets.target
